# gh-runner-metrics.service
[Unit]
Description=GitHub Actions Runner Metrics Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/bv/gh-runner-scaler/metrics.py
Environment=GRAFANA_CLOUD_API_KEY=your-grafana-cloud-api-key-here
Environment=LOKI_PUSH_URL=your-loki-push-url-here
Environment=LOKI_USERNAME=your-loki-username-here
User=bv
